# api/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi --no-site --ini /etc/uwsgi/api.ini
